# check-metered-connection.service
[Unit]
Description=Metered connection status checker

[Service]
Type=oneshot
ExecStart=check-metered-connection.sh
