# unmetered-connection.target
[Unit]
Description=Unmetered connection
Requires=check-metered-connection.service
After=check-metered-connection.service
