# unmetered-connection.target
[Unit]
Description=Unmetered connection
Requires=check-metered-connection-status.service
After=check-metered-connection-status.service
